# upstartstreaming.conf
description "starts a script for streaming server"

start on filesystem and started networking
stop on runlevel [!2345]

script
	/usr/local/bin/initializestreaming.sh
end script
